# aws/hackaton.service
[Unit]
Description=Proyecto para el hackaton

[Service]
ExecStart=/usr/bin/npm start --prefix

[Install]
WantedBy=mongod.target
